# deploy/systemd/tradingbot.service
[Unit]
# Replace __PROJECT_ROOT__ before loading this unit.
Description=TradingBot automated options trading daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_ROOT__
EnvironmentFile=__PROJECT_ROOT__/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=__PROJECT_ROOT__/.venv/bin/python __PROJECT_ROOT__/main.py --live --yes
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
