# deploy/tradingbot.service
[Unit]
Description=TradingBot — Automated Options Trading
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tradingbot
Group=tradingbot
WorkingDirectory=/opt/tradingbot
ExecStart=/opt/tradingbot/venv/bin/python3 main.py run paper
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

# Security hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/tradingbot

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
